# docs/guides/agents/systemd/templates/com.tokdash.tokdash.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
  <dict>
    <key>Label</key>
    <string>com.tokdash.tokdash</string>

    <!--
      NOTE:
      For best reliability, replace `/usr/bin/env` + `tokdash` with the absolute path
      to your tokdash binary (find it with `which tokdash`).
    -->
    <key>ProgramArguments</key>
    <array>
      <string>/usr/bin/env</string>
      <string>tokdash</string>
      <string>serve</string>
      <string>--bind</string>
      <string>127.0.0.1</string>
      <string>--port</string>
      <string>55423</string>
      <!-- Headless service: don't try to launch a browser. -->
      <string>--no-open</string>
    </array>

    <key>EnvironmentVariables</key>
    <dict>
      <key>PATH</key>
      <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
      <key>PYTHONUNBUFFERED</key>
      <string>1</string>
    </dict>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <!-- launchd requires absolute paths; /tmp is a safe default. -->
    <key>StandardOutPath</key>
    <string>/tmp/tokdash.out.log</string>
    <key>StandardErrorPath</key>
    <string>/tmp/tokdash.err.log</string>
  </dict>
</plist>

